# vapor_server.ini
; vapor_server.ini
[program:PROGRAM_NAME]
command=vapor run serve --hostname=127.0.0.1 --port=8080
autostart=true
autorestart=true
startretries=5
user=root

[supervisord]
